# supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:cowsay]
command=/var/script/cowsay.sh

[program:apache2]
command=httpd-foreground
